# .docker/php7-fpm/supervisor-hubspot.conf
[program:wicked-report-hubspot-listener]
command=php ./bin/console app:beanstalk:queue-hubspot-fetcher
numprocs=1
autostart=true
autorestart=true
directory=/var/www/symfony
